# linux-mastery-exercises/tier-2-operations/exercise-2.6-overprivileged-agent/insecure-server/etc/systemd/system/report-generator.service
[Unit]
Description=Report Generator Agent
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/agents/report-generator/src/main.py
EnvironmentFile=/opt/agents/report-generator/config/.env
User=root
# SECURITY ISSUE: Running as root!
WorkingDirectory=/opt/agents/report-generator
Restart=always

[Install]
WantedBy=multi-user.target
