# linux-mastery-exercises/tier-2-operations/exercise-2.5-lock-it-down/insecure-deployment/current-unit-file.service
[Unit]
Description=Social Media Monitor
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/agents/social-monitor/src/main.py
EnvironmentFile=/opt/agents/social-monitor/config/.env
User=root
# SECURITY ISSUE: Running as root!
WorkingDirectory=/opt/agents/social-monitor
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
